# xremap/xremap.service
[Unit]
Description=xremap

[Service]
Restart=always
ExecStart=xremap %h/linuxc/xremap/config.yaml

[Install]
WantedBy=default.target